# conf/supervisord.conf
[supervisord]
nodaemon=true
loglevel=error

[program:spacewalk-setup]
command=/opt/spacewalk.sh
autorestart=false
